# deploy/rz-admin.service
[Unit]
Description=RustZen Admin
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=/opt/rz/config/rz.env
ExecStart=/opt/rz/bin/rz admin serve
WorkingDirectory=/opt/rz
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
