# deploy/rz-reports.service
[Unit]
Description=RustZen Reports Worker
After=network.target rz-recovery.service
PartOf=rz.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=/opt/rz/config/rz.env
ExecCondition=/usr/bin/test ! -e /opt/rz/data/recovery-blocked
ExecStart=/opt/rz/current/bin/rz-reports serve
WorkingDirectory=/opt/rz
Restart=on-failure
RestartSec=5

[Install]
WantedBy=rz.target
