# deploy/rz-recovery.service
[Unit]
Description=RustZen Interrupted Update Recovery
After=network.target
Before=rz-admin.service rz-monitor.service rz-insights.service rz-reports.service
PartOf=rz.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
EnvironmentFile=/opt/rz/config/rz.env
ExecStart=/opt/rz/current/bin/rz-admin update recover
WorkingDirectory=/opt/rz
Restart=on-failure
RestartSec=2

[Install]
WantedBy=rz.target
